# meta-composeos/recipes-core/composeos-openrc/files/composeos-run.initd
#!/sbin/openrc-run

description="ComposeOS compose service"

depend() {

  after chrony-wait

}

start()
{
  ebegin "composeOs compose service stacks..."
  /lib/composeos/runrc.sh
  eend $?
}
